# .config/containers/systemd/paperless-tika.container
[Unit]
Description=tika for Paperless NGX

[Container]
Pod=paperlessNGX.pod
Image=ghcr.io/paperless-ngx/tika:latest
